# overlay-inaugural/chromeos-base/chromeos-bsp-inaugural/files/wifi_init/99-start-wifi.rules
KERNEL=="wlan0", ACTION=="add", RUN+="/lib/udev/start_wifi.sh"
