# baseboard-inaugural/chromeos-base/autoexpand-user-partition/files/autoexpand-user-partition.conf
# Copyright (c) 2022 The FydeOS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Job to expand user's partition"
author          "yang@fydeos.io"

start on started boot-services

task
oom score 0

script
  user_partition_mount_point="/mnt/stateful_partition"
  expanded_file=${user_partition_mount_point}/.user_partition_expanded
  user_partition=$(/bin/mount |grep ' /mnt/stateful_partition ')
  if [ ! -f $expanded_file ]; then
    /sbin/resize2fs ${user_partition%% on*}
    /usr/bin/touch $expanded_file
  fi
end script
