# overlay-inaugural/chromeos-base/chromeos-bsp-inaugural/files/ethernet_issue_workaround.conf
# Copyright (c) 2022 The FydeOS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Workaround before fixing shill's issus about ethernets"
author          "yang@fydeos.io"

# VPD (Vital Product Data) information is shown on chrome://system for
# user, feedback, and partner's customer service center.
start on started boot-services
task
oom score 0

script
  logger -t "${UPSTART_JOB}" "Bring up ethernets"
  /bin/ifconfig eth0 up
  /bin/ifconfig eth1 up
end script
